# board/cryocon/m54/rootfs_overlay/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
null::sysinit:/bin/mount -t proc proc /proc
#null::sysinit:/bin/mount -o remount,rw / # REMOUNT_ROOTFS_RW
#null::sysinit:/bin/mount -t devtmpfs dev /dev
null::sysinit:/bin/mkdir -p /dev/pts
null::sysinit:/bin/mkdir -p /dev/shm
::sysinit:/sbin/fsck -p /dev/mmcblk0p3
null::sysinit:/bin/mount -a
::sysinit:/usr/bin/head -1 /sys/class/graphics/fb0/modes > /sys/class/graphics/fb0/mode
::sysinit:/usr/bin/psplash -n &
::sysinit:/bin/echo 6 > /sys/class/backlight/backlight/brightness
null::sysinit:/bin/hostname -F /etc/hostname
::once:/sbin/modprobe ft5x06-ts

# now run any rc scripts
::once:/etc/init.d/rcS

::sysinit:/sbin/modprobe xilinx-sscu
::sysinit:/usr/bin/load_fpga /boot/fpga1.bit /dev/fpga1
::sysinit:/sbin/modprobe -r xilinx-sscu
::sysinit:/sbin/modprobe rtc-at91rm9200
::sysinit:/sbin/modprobe spi-atmel

::sysinit:/bin/mkdir -p /tmp/.X11-unix
::respawn:env DISPLAY=:0 /usr/bin/cryocond
::respawn:/usr/bin/Xorg -background none
::respawn:/bin/su -s /bin/sh -c "/usr/bin/gunicorn -k flask_sockets.worker -b localhost\:8000 cryocond:app" - www-data
::respawn:/usr/sbin/crond -l 1 -f
::respawn:/usr/bin/LinkMD

# Put a getty on the serial port
# ttyS0::respawn:/bin/sh
ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100 # GENERIC_SERIAL

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
null::shutdown:/sbin/modprobe at91-reset
null::shutdown:/etc/init.d/rcK
null::shutdown:/usr/bin/killall klogd
null::shutdown:/usr/bin/killall syslogd
null::shutdown:/bin/umount -a -r
null::shutdown:/sbin/swapoff -a

