# board/cryocon/m18_g35_ek/overlay/lib/udev/rules.d/85-usbmount.rules
# Rules for USBmount -*- conf -*-

ACTION=="add",SUBSYSTEM=="block",KERNEL=="sd[a-h]1",	RUN+="usb_mount /dev/%k"
ACTION=="remove",SUBSYSTEM=="block",KERNEL=="sd[a-h]1",	RUN+="usb_unmount /dev/%k"
