# board/cryocon/m18_g35_ek/overlay/etc/network/if-pre-up.d/firewall
iptables-restore << EOF
*filter
:INPUT DROP
-A INPUT -p udp -j ACCEPT
-A INPUT -p tcp -j ACCEPT
COMMIT
EOF
