# board/cryocon/m18_g35_ek/overlay/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
null::sysinit:/bin/mount -t proc proc /proc
#null::sysinit:/bin/mount -o remount,rw / # REMOUNT_ROOTFS_RW
#null::sysinit:/bin/mount -t devtmpfs dev /dev
null::sysinit:/bin/mkdir -p /dev/pts
null::sysinit:/bin/mkdir -p /dev/shm
::sysinit:/sbin/fsck -p /dev/mmcblk0p3
null::sysinit:/bin/mount -a
null::sysinit:/bin/hostname -F /etc/hostname
::sysinit:/bin/mknod /dev/spidev0.0 c 153 0
::sysinit:/bin/mknod /dev/spidev1.0 c 153 1

# now run any rc scripts
::once:/etc/init.d/rcS

::respawn:/usr/bin/cryocond
::respawn:/bin/sh -c "cd /var && /usr/bin/cryocon-fcgi.py"
::respawn:/usr/sbin/crond -l 1 -f

# Put a getty on the serial port
# ttyS0::respawn:/bin/sh
ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100 # GENERIC_SERIAL

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
null::shutdown:/etc/init.d/rcK
null::shutdown:/usr/bin/killall klogd
null::shutdown:/usr/bin/killall syslogd
null::shutdown:/bin/umount -a -r
null::shutdown:/sbin/swapoff -a

