# board/cryocon/common/rootfs_overlay/etc/network/if-up.d/default-route
#!/bin/sh

ip route add default dev eth0 && echo "Howdy" || echo "Oops"
